# automation/com.tradingagent.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.deltaforge.daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PROJECT_ROOT__/venv/bin/python3</string>
        <string>__PROJECT_ROOT__/main.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PROJECT_ROOT__</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>25</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>__PROJECT_ROOT__/logs/launchd_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__PROJECT_ROOT__/logs/launchd_stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__PROJECT_ROOT__/venv/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>TZ</key>
        <string>Asia/Kolkata</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
